# scripts/juno@.service
[Unit]
Description=Juno public-agent runtime (%i)
Documentation=https://github.com/idchain-world/juno
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=juno
Group=juno
WorkingDirectory=/opt/juno
EnvironmentFile=/etc/juno/%i.env
ExecStart=/usr/bin/node /opt/juno/dist/server.js
Restart=on-failure
RestartSec=5

# Per-instance resource limits. One misbehaving juno cannot starve the VPS.
# Tune these via drop-ins (/etc/systemd/system/juno@<name>.service.d/) if
# a specific instance needs more headroom. Defaults sized for CX22 (4 GB / 2 vCPU)
# carrying ~3-5 junos.
MemoryMax=512M
MemoryHigh=384M
CPUWeight=100
TasksMax=256
LimitNOFILE=4096

# Journal log volume bound — one chatty instance should not fill /var/log.
LogRateLimitIntervalSec=10s
LogRateLimitBurst=2000

# Sandboxing — Ubuntu 24.04 systemd supports all of these out of the box.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictRealtime=true
SystemCallArchitectures=native
ReadWritePaths=/opt/juno/instances/%i

[Install]
WantedBy=multi-user.target
